# pi/systemd/aquaponics-uploader.service
[Unit]
Description=Aquaponics Pi - SQLite to cloud uploader
After=network.target aquaponics-app.service
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/smart-aquaponics/pi
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/pi/smart-aquaponics/pi/uploader.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
